# examples/upstart/celeryd-cron.conf
description "FuzzManager cron celeryd"
start on runlevel [2345]
stop on runlevel [06]

setuid ubuntu
setgid ubuntu

chdir /path/to/FuzzManager/server

respawn

exec celery -A celeryconf -l INFO -c 1 -n cron@%h -Q cron worker
